# config/upstart/mtlfe-server-8003.conf
#!upstart
description "mtlfe server 03"
author      "quickredfox"

start on startup
stop on shutdown

respawn

exec sudo -iu mtlfe /usr/bin/env node /var/mtlfe/dispatch/lib/server.js -h mltfe.org -p 8003 >> /var/log/mtlfe-dispatch.log 2>&1 
