# config/upstart/mtlfe-hook-9001.conf
#!upstart
description "mtlfe server 03"
author      "quickredfox"

start on startup
stop on shutdown

respawn

exec sudo -iu mtlfe /usr/bin/env node /var/mtlfe/dispatch/lib/hook.js -h mltfe.org -p 9001 >> /var/log/mtlfe-hook.log 2>&1 
